# eero-bypass.service.default
[Unit]
Description=delphijustin EERO Bypass
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/eero-netchecker.sh nowrite
Restart=on-failure
User=user_name
Group=user_name

[Install]
WantedBy=multi-user.target
